# nginx/configs/etc/init/ecoservice.conf
description "Eco benefit service"

start on vagrant-mounted
stop on shutdown

respawn
respawn limit 50 5

kill timeout 5

chdir /usr/local/ecoservice/ecoservice

exec ./ecobenefits 2>&1 1>/var/log/ecoservice
